# debian/preinst
#!/bin/sh
# preinst script for deepamehta
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <new-preinst> `install'
#        * <new-preinst> `install' <old-version>
#        * <new-preinst> `upgrade' <old-version>
#        * <old-preinst> `abort-upgrade' <new-version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

BACKUPDATE=$( date +%F-%H:%M:%S )
DMVERSION='@@VERSION@@'
PREVIOUSVERSION=$( find /var/lib/deepamehta/dmx-db/ -name "*.dm-version.*" | cut -d'.' -f'3-' | sort -V | tail -n1 )

if [ -f /var/run/deepamehta/deepamehta.pid ]; then
    invoke-rc.d deepamehta stop
fi

case "$1" in
    install)
    useradd --system --user-group --create-home -d'/var/lib/deepamehta' -s'/bin/false' -c'DeepaMehta - Cope With Complexity' deepamehta
    if [ -d /var/lib/deepamehta/dmx-db ]; then
	    if [ ! -f /var/lib/deepamehta/dmx-db/.dm-version.$DMVERSION ]; then
		    mv /var/lib/deepamehta/dmx-db /var/lib/deepamehta/dmx-db-$PREVIOUSVERSION-backup-$BACKUPDATE
		    mkdir /var/lib/deepamehta/dmx-db
		    chown deepamehta:deepamehta /var/lib/deepamehta/dmx-db
		    touch /var/lib/deepamehta/dmx-db/.dm-version.$DMVERSION
		    echo "Your existing database has been moved to /var/lib/deepamehta/dmx-db-$PREVIOUSVERSION-backup-$BACKUPDATE."
	    fi
    else
	    mkdir /var/lib/deepamehta/dmx-db
	    chown deepamehta:deepamehta /var/lib/deepamehta/dmx-db
	    touch /var/lib/deepamehta/dmx-db/.dm-version.$DMVERSION
    fi
    if [ -d /var/cache/deepamehta/bundle-cache ]; then
	    rm -rf /var/cache/deepamehta/bundle-cache
    fi
    if [ -d /var/cache/deepamehta/deepamehta-bundle-cache ]; then
	    rm -rf /var/cache/deepamehta/deepamehta-bundle-cache
    fi
    ;;

    upgrade)
    # existing DB is < DM 4.1
    if [ -d /var/lib/deepamehta/dmx-db ] && [ ! -f /var/lib/deepamehta/dmx-db/.dm-version.$DMVERSION ] && [ -f /var/lib/deepamehta/dmx-db/.dm-version.4.0* ]; then
	    mv /var/lib/deepamehta/dmx-db /var/lib/deepamehta/dmx-db-$PREVIOUSVERSION-backup-$BACKUPDATE
	    mkdir /var/lib/deepamehta/dmx-db
	    chown deepamehta:deepamehta /var/lib/deepamehta/dmx-db
	    touch /var/lib/deepamehta/dmx-db/.dm-version.$DMVERSION
	    echo "Your existing database has been moved to /var/lib/deepamehta/dmx-db-$PREVIOUSVERSION-backup-$BACKUPDATE."
	    rm /var/lib/deepamehta/dmx-db/.dm-version.$PREVIOUSVERSION
    # existing DB is >= DM 4.1
    elif [ -d /var/lib/deepamehta/dmx-db ] && [ ! -f /var/lib/deepamehta/dmx-db/.dm-version.$DMVERSION ] && [ -f /var/lib/deepamehta/dmx-db/.dm-version.4.1* ]; then
        cp -a /var/lib/deepamehta/dmx-db /var/lib/deepamehta/dmx-db-$PREVIOUSVERSION-backup-$BACKUPDATE
        touch /var/lib/deepamehta/dmx-db/.dm-version.$DMVERSION
		echo "Your existing database has been copied to /var/lib/deepamehta/deepamehta-db-$PREVIOUSVERSION-backup-$BACKUPDATE."
		rm /var/lib/deepamehta/dmx-db/.dm-version.$PREVIOUSVERSION
    fi
    if [ -d /var/cache/deepamehta/bundle-cache ]; then
	    rm -rf /var/cache/deepamehta/bundle-cache
    fi
    if [ -d /var/cache/deepamehta/deepamehta-bundle-cache ]; then
	    rm -rf /var/cache/deepamehta/deepamehta-bundle-cache
    fi
    ;;

    abort-upgrade)
    ;;

    *)
        echo "preinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
